# _today.conf
[program:today]
command = /home/projects/today/_deploy.sh
